# telescope-local.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN"
    "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.galmungral.telescope-local</string>
    <key>ProgramArguments</key>
    <array>             
        <string>/usr/local/bin/telescope-local</string>
        <string>3030</string>
        <string>telescope.remote</string>
        <string>3030</string>
    </array>
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>