# build/1.1.0/telescope.service
[Service]
ExecStart=/usr/local/bin/telescope-remote 3030
